# services/nadeko/run
#!/bin/bash
umask 0002
abc mono /root/NadekoBotNadekoBot.exe
