# services/nadeko/run
#!/bin/bash
umask 0002
mono /root/NadekoBotNadekoBot.exe
