# services/nadeko/run
#!/bin/bash
umask 0002
cd /root/NadekoBot
mono NadekoBot.exe