# services/nadeko/run
#!/bin/bash
umask 0002
cd /root/NadekoBot
XDG_CONFIG_HOME="/config/xdg" exec /sbin/setuser abc mono NadekoBot.exe
